# deploy/carolyn-merrick-api.service
# /etc/systemd/system/carolyn-merrick-api.service  (on the App CT)
[Unit]
Description=Carolyn & Merrick wedding API (Fastify)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/carolyn-merrick-api
EnvironmentFile=/opt/carolyn-merrick-api/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
User=cmapi
Group=cmapi
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/carolyn-merrick-api/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
